# etc/hotplug.d/iface/90-miniupnpd_updater
#!/bin/bash
if [ "$ACTION" = "ifup" ] && [ "$INTERFACE" = "wan" ]; then
    /etc/scripts/set_upnp_public_ip.sh
fi
